# services/gitea/supervisor.conf
[program:gitea]
# directory=/usr/local/bin/
command=/usr/local/bin/gitea web
	--work-path=<BASE_DIR>/
	--custom-path=<BASE_DIR>/custom/
	-c <BASE_DIR>/custom/conf/app.ini
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/gitea/stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/gitea/stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
user = git
environment = HOME="/home/git", USER="git"
